# supervisord.conf
[program:thyme]
command=node ./lib/server
directory=/srv/thyme
