# supervisord.conf
[program:thyme]
command=node ./lib/server
environment=MONGO_HOST=localhost,MONGO_DATABASE=thyme,RABBITMQ_HOST=localhost
directory=/srv/thyme
