# supervisord.conf
[program:master]
command=node ./lib/server
directory=/srv/thyme
